# piwik/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "service apache2 start"

[program:mariadb]
command=/bin/bash -c "service mysql start"
